# app/base-files/files/etc/init.d/firewall
#!/bin/sh /etc/rc.common
# START=25

start
{
	iptables -t nat -A POSTROUTING -o eth1 -j MASQUERADE
}

