# app/base-files/files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
ttyMSM0::askfirst:-/bin/sh
::shutdown:/etc/init.d/rcS K shutdown
